# rootfs/etc/s6-overlay/s6-rc.d/check-playback-duration/finish
#!/usr/bin/with-contenv bash

# dont restart this service
exit 125