# rootfs/etc/s6-overlay/s6-rc.d/spotifyd/finish
#!/usr/bin/with-contenv bash

s6-permafailon 60 5 0-255